# crontab/homebot
# $Id$
0    4   * * *  /usr/local/bin/krenew
*/5  *   * * *  /usr/local/bin/tursom_checkin
0    */8 * * *  /usr/local/bin/other_checkin
