# crontab/homebot
# $Id$
0    4 * * *  /usr/local/bin/krenew
*/5  * * * *  /usr/local/bin/tursom_checkin
*/20 * * * *  /usr/local/bin/update_demonoid
