# crontab/tmaher.ps24007
55 23 * * * $HOME/garden/expo/eat_tweets.sh
